# backup-services/service/retroarch-backup.timer
[Unit]
Description=Run retroarch backup every 30 minutes

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
