# backup-services/service/yuzu-backup.timer
[Unit]
Description=Run yuzu backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
